# scripts/docker/supervisord/tca_service.conf
[group:tca-server]
programs = main-server, main-worker, main-beat, analysis-server, analysis-worker, file-server, login-server, scmproxy, nginx


[program:main-server] 
command=bash -c 'sleep 5 && gunicorn codedog.wsgi -c main.gunicorn.conf.py'
directory=/CodeAnalysis/server/projects/main
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/CodeAnalysis/server/projects/main/log/supervisor_main_server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/CodeAnalysis/server/projects/main/log/supervisor_main_server_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:main-worker] 
command=bash -c 'sleep 5 && celery -A codedog worker -l DEBUG'
directory=/CodeAnalysis/server/projects/main
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/CodeAnalysis/server/projects/main/log/supervisor_main_worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/CodeAnalysis/server/projects/main/log/supervisor_main_worker_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:main-beat] 
command=bash -c 'sleep 5 && celery -A codedog beat -S redbeat.RedBeatScheduler -l DEBUG'
directory=/CodeAnalysis/server/projects/main
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/CodeAnalysis/server/projects/main/log/supervisor_main_beat.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/CodeAnalysis/server/projects/main/log/supervisor_main_beat_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:analysis-server] 
command=bash -c 'sleep 5 && gunicorn codedog.wsgi -c analysis.gunicorn.conf.py'
directory=/CodeAnalysis/server/projects/analysis
priority=10
startsecs=10
stopwaitsecs=1
autostart=true          
autorestart=true
stopasgroup=true
stdout_logfile=/CodeAnalysis/server/projects/analysis/log/supervisor_analysis_server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/CodeAnalysis/server/projects/main/log/supervisor_analysis_server_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:analysis-worker] 
command=bash -c 'sleep 5 && celery -A codedog worker -l DEBUG'
directory=/CodeAnalysis/server/projects/analysis
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/CodeAnalysis/server/projects/analysis/log/supervisor_analysis_worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/CodeAnalysis/server/projects/main/log/supervisor_analysis_worker_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:file-server] 
command=bash -c 'sleep 5 && gunicorn codedog_file_server.wsgi -c file.gunicorn.conf.py'
directory=/CodeAnalysis/server/projects/file
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/CodeAnalysis/server/projects/file/log/supervisor_file_server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/CodeAnalysis/server/projects/main/log/supervisor_file_server_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:login-server] 
command=bash -c 'sleep 5 && gunicorn apps.wsgi -c login.gunicorn.conf.py'
directory=/CodeAnalysis/server/projects/login
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/CodeAnalysis/server/projects/login/log/supervisor_login_server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/CodeAnalysis/server/projects/main/log/supervisor_login_server_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:scmproxy] 
command=python proxyserver.py
directory=/CodeAnalysis/server/projects/scmproxy
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
exitcodes=0                                          
stopsignal=QUIT                                      
stopwaitsecs=1    
stdout_logfile=/CodeAnalysis/server/projects/scmproxy/logs/supervisor_scmproxy.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/CodeAnalysis/server/projects/main/log/supervisor_scmproxy_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10                             


[program:nginx]
command=bash -c "sleep 5 && /usr/sbin/nginx -g 'daemon off;'"
priority=20
startsecs=10
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/nginx/supervisor_access.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/nginx/supervisor_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10    
