# scripts/docker/supervisord/tca_service.conf
[group:tca-server]
programs = main-server, main-worker, main-beat, analysis-server, analysis-worker, file-server, login-server, scmproxy, nginx


[program:main-server] 
command=bash -c 'sleep 5 && gunicorn codedog.wsgi -c main.gunicorn.conf.py'
directory=/CodeAnalysis/server/projects/main
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/var/log/tca/servers/main/server_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/tca/servers/main/server_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10

[program:main-worker] 
command=bash -c 'sleep 5 && celery -A codedog worker -l DEBUG'
directory=/CodeAnalysis/server/projects/main
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/var/log/tca/servers/main/worker_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/tca/servers/main/worker_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:main-beat] 
command=bash -c 'sleep 5 && celery -A codedog beat -S redbeat.RedBeatScheduler -l DEBUG'
directory=/CodeAnalysis/server/projects/main
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/var/log/tca/servers/main/beat_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/tca/servers/main/beat_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:analysis-server] 
command=bash -c 'sleep 5 && gunicorn codedog.wsgi -c analysis.gunicorn.conf.py'
directory=/CodeAnalysis/server/projects/analysis
priority=10
startsecs=10
stopwaitsecs=1
autostart=true          
autorestart=true
stopasgroup=true
stdout_logfile=/var/log/tca/servers/analysis/server_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/tca/servers/analysis/server_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:analysis-worker] 
command=bash -c 'sleep 5 && celery -A codedog worker -l DEBUG'
directory=/CodeAnalysis/server/projects/analysis
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/var/log/tca/servers/analysis/worker_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/tca/servers/analysis/worker_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:file-server] 
command=bash -c 'sleep 5 && gunicorn codedog_file_server.wsgi -c file.gunicorn.conf.py'
directory=/CodeAnalysis/server/projects/file
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/var/log/tca/servers/file/server_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/tca/servers/file/server_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:login-server] 
command=bash -c 'sleep 5 && gunicorn apps.wsgi -c login.gunicorn.conf.py'
directory=/CodeAnalysis/server/projects/login
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
stopasgroup=true
stdout_logfile=/var/log/tca/servers/login/server_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/tca/servers/login/server_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10


[program:scmproxy] 
command=python proxyserver.py
directory=/CodeAnalysis/server/projects/scmproxy
priority=10
startsecs=10
stopwaitsecs=1
autostart=true
autorestart=true
exitcodes=0                                          
stopsignal=QUIT                                      
stopwaitsecs=1    
stdout_logfile=/var/log/tca/servers/scmproxy/server_stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/tca/servers/scmproxy/server_stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10                             


[program:nginx]
command=bash -c "sleep 5 && /usr/sbin/nginx -g 'daemon off;'"
priority=20
startsecs=10
autostart=true
autorestart=true
stdout_logfile=/var/log/tca/nginx/access.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/tca/nginx/error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10    


[program:client]
command=bash -c "sleep 10 && source ../scripts/config.sh && python3 codepuppy.py start -t $CODEDOG_TOKEN"
directory=/CodeAnalysis/client/
priority=30
startsecs=10
autostart=true
autorestart=true
stdout_logfile=/var/log/tca/client/client_start.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stderr_logfile=/var/log/tca/client/client_error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10    
